# install/pi-kiosk.conf
description "Kiosk Playlist"
author "Tom Jenkinson"

respawn
start on net-device-up
stop on shutdown

script
	export HOME="/home/pi"
	exec /home/pi/app/start.sh
end script