# install/pi-kiosk.conf
description "Kiosk Playlist"
author "Tom Jenkinson"

respawn
setuid pi
start on net-device-up
stop on shutdown

exec /bin/bash -c '/home/pi/app/start.sh'
